# services.d/opencode/run
#!/command/with-contenv bash

set -e

echo "Starting opencode..."

# ensure opencode config exists; if missing, write a sensible default
CONFIG_PATH="$HOME/.config/opencode/opencode.json"
if [ ! -f "${CONFIG_PATH}" ]; then
	echo "opencode config not found; writing default to ${CONFIG_PATH}"
	mkdir -p "$(dirname "${CONFIG_PATH}")"
	cp $HOME/default-config/opencode.json "${CONFIG_PATH}"
fi

exec opencode web --port "${OPENCODE_PORT:-4096}" --hostname "${OPENCODE_HOSTNAME:-0.0.0.0}"
